# daemons/update-services/update_services.service
[Unit]
Description=Update machine services

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/local/bin/update_services.sh
