# daemons/network-watchdog/network-watchdog.service
[Unit]
Description=Restart network service after wifi issues

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/local/bin/network-watchdog.sh
